# deploy/natalisbot.service
[Unit]
Description=NatalisBot - Telegram Bot для юридических услуг
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/var/telegramBots/NatalisBot
Environment="PATH=/var/telegramBots/NatalisBot/venv/bin"
ExecStart=/var/telegramBots/NatalisBot/venv/bin/python bot.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target

